# server/simjiang-reminder/simjiang-reminder.service
[Unit]
Description=SimJiang E2EE Cloud Reminder Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/simjiang-reminder
Environment=SIMJ_BASE=/opt/simjiang-reminder
Environment=SIMJ_HOST=0.0.0.0
Environment=SIMJ_PORT=8787
ExecStart=/usr/bin/python3 /opt/simjiang-reminder/server.py
Restart=always
RestartSec=3
# Fresh schema lives in data.db; legacy files are ignored

[Install]
WantedBy=multi-user.target
